# tests/templates/wwwkindos.conf
[program:{USERNAME}]
command=bash /home/{USERNAME}/gunicorn_start
user={USERNAME}
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/{USERNAME}/gunicorn-error.log
